# tlp.upstart.in
# tlp - system startup/shutdown
#
# Copyright (c) 2016 Thomas Koch <linrunner at gmx.net>
# This software is licensed under the GPL v2 or later.

description "tlp"

start on ( virtual-filesystems and runlevel [2345] )
stop on runlevel [!2345]

env TLP=@TLP_SBIN@/tlp

pre-start script
    [ -x $TLP ] || exit 4
    $TLP init start

end script

post-stop script
    [ -x $TLP ] || exit 4
    $TLP init stop

end script


