# tlp.rules
# tlp - handle added usb devices
ACTION=="add", SUBSYSTEM=="usb", RUN+="/lib/udev/tlp-udev %k"
