# tlp.service.in
# tlp - systemd startup/shutdown service
#
# Copyright (c) 2017 Thomas Koch <linrunner at gmx.net>
# This software is licensed under the GPL v2 or later.

[Unit]
Description=TLP system startup/shutdown
Wants=bluetooth.service NetworkManager.service
After=multi-user.target bluetooth.service NetworkManager.service
Before=shutdown.target
Documentation=http://linrunner.de/tlp

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@TLP_SBIN@/tlp init start
ExecStop=@TLP_SBIN@/tlp init stop

[Install]
WantedBy=multi-user.target
