# debian/tlp.postinst
#!/bin/sh
set -e

#DEBHELPER#

SERVICES2MASK="systemd-rfkill.service systemd-rfkill.socket"

case "$1" in
    configure)
        # Mask conflicting systemd services
        for svc in $SERVICES2MASK; do
            deb-systemd-helper mask $svc 2> /dev/null || true
        done
        ;;
esac
